# etc/nightwatch-backup/schedules/every-6-hours.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

0 */6 * * * root /usr/local/bin/nightwatch-backup run >/var/log/nightwatch-backup/cron.log 2>&1
