# install/templates/systemd/nightwatch-backup.timer
[Unit]
Description=Nightwatch Backup Scheduled Backups

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target